# boot/scripts/backup-apkovl.timer
[Unit]
Description=Backup APKOVL's weekly

[Timer]
OnCalendar=Wed *-*-* 07:00
Persistent=true

[Install]
WantedBy=timers.target
